# roles/slurm_master/templates/slurmrestd.service.j2
# =============================================================================
# SLURMRESTD SYSTEMD SERVICE FILE - GENERATED BY ANSIBLE  
# =============================================================================

[Unit]
Description=Slurm REST API Daemon
Documentation=man:slurmrestd(8)
After=network.target munge.service slurmctld.service
Wants=network.target
Requires=munge.service slurmctld.service

[Service]
Type=notify
EnvironmentFile=-/etc/default/slurmrestd
User={{ slurm_restd_user | default('slurmrestd') }}
Group={{ slurm_group | default('slurm') }}

ExecStart={{ slurm_install_prefix }}/sbin/slurmrestd $SLURMRESTD_OPTIONS {{ slurm_master_address }}:{{ slurm_slurmrestd_port }}
ExecReload=/bin/kill -HUP $MAINPID

# Restart behavior
Restart=on-failure
RestartSec=10s
TimeoutStopSec=30s

# Security settings
NoNewPrivileges=yes
PrivateTmp=yes
ProtectHome=yes
ProtectSystem=strict

# Make log directory writable
ReadWritePaths=/var/log/slurmrestd
{% if slurm_restd_config_dir is defined %}
ReadWritePaths={{ slurm_restd_config_dir }}
{% endif %}

# Resource limits
LimitNOFILE=65536
LimitNPROC=65536

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=slurmrestd

[Install]
WantedBy=multi-user.target